# pterocleaner.service
[Unit]
Description=PteroCleaner - Orphaned backup cleaner for Pterodactyl
After=network.target

[Service]
Type=simple
WorkingDirectory=/etc/pterocleaner
ExecStart=/etc/pterocleaner/venv/bin/python /etc/pterocleaner/pterocleaner.py
Restart=always
RestartSec=5
User=root

[Install]
WantedBy=multi-user.target
